# scripts/jobd-broker.service
[Unit]
Description=jobd broker (docker compose stack)
After=docker.service tailscaled.service network-online.target
Requires=docker.service
Wants=tailscaled.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
# EDIT: set to wherever you cloned this repo on the broker host.
WorkingDirectory=/opt/jobd
# Wait until the Tailscale IP is actually programmed onto an interface.
# tailscaled.service being "active" only means the daemon is running —
# the overlay address is assigned asynchronously after handshake. Without
# this gate, `compose up -d` races and fails with
# "cannot assign requested address" because the bind in docker-compose.yml
# targets that IP. 60 iterations * 2s = 120s ceiling.
#
# EDIT: replace 100.64.0.1 with this host's tailscale IP (the value you set
# for JOBD_HOST). If you run a single-host loopback broker (JOBD_HOST=127.0.0.1),
# delete this ExecStartPre line entirely — loopback is always present.
ExecStartPre=/bin/bash -c 'for i in $(seq 1 60); do ip -4 -o addr show | grep -q "100.64.0.1" && exit 0; sleep 2; done; echo "tailscale IP 100.64.0.1 never appeared after 120s" >&2; exit 1'
# Defensive cleanup: a previously-failed boot can leave a stale "Created"
# container that compose up -d treats as already-present and refuses to
# fix port bindings on. Leading dash makes this best-effort.
ExecStartPre=-/usr/bin/docker compose down --remove-orphans
# --build is load-bearing, not belt-and-braces. docker-compose.yml builds from
# source (`build: .`) and tags the result `jobd:latest` with no registry behind
# it, so `up -d` alone happily reuses whatever `jobd:latest` already exists —
# i.e. it restarts the OLD code after a `git pull`. Without --build, a boot
# following an upgrade silently serves the previous version, and `job ping`
# reports a version that doesn't match the checkout. Rebuilding is a no-op
# (layer cache) when nothing changed.
ExecStart=/usr/bin/docker compose up -d --build --remove-orphans
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=180
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
